# detections/endpoint/cronjob-suspicious-path.yaml
id: det-endpoint-023
name: Cron Job Created Pointing to /tmp or /dev/shm
description: AiSOC v1 curated detection. Triggers on the endpoint signal described by 'Cron Job Created
  Pointing to /tmp or /dev/shm'. Watch the 1 documented false-positive case before tuning.
version: 1.0.0
severity: high
tags:
- mitre.attack.t1053.003
- tlp.white
category: endpoint
log_source:
  product: linux
  service: auditd
detection:
  condition: |
    event_type == "PATH"
    AND file_path STARTSWITH_ANY ["/etc/cron", "/var/spool/cron"]
    AND command_text CONTAINS_ANY ["/tmp/", "/dev/shm/", "curl", "wget"]
false_positives:
- Vendor packaged cron jobs (must be from package manager)
playbook: tpl-persistence
enabled: true
author: AiSOC
created: '2026-05-03'
modified: '2026-05-03'
