# detections/endpoint/fim-sensitive-file-modified.yaml
id: det-endpoint-281
name: Sensitive System File Modified (FIM)
description: >
  A file in a sensitive system path (/etc/passwd, /etc/shadow, /etc/sudoers,
  /etc/crontab, /etc/hosts, /etc/ssh/sshd_config, /etc/ld.so.*)  was created,
  deleted, or modified as reported by osquery file_events via the AiSOC FIM
  pipeline.  Adversaries commonly target these files for credential theft
  (T1003.008), persistence (T1037, T1053.003), lateral movement, and privilege
  escalation (T1548.001).
version: 1.0.0
severity: high
tags:
  - mitre.attack.t1003.008   # OS Credential Dumping: /etc/passwd and /etc/shadow
  - mitre.attack.t1037        # Boot or Logon Initialization Scripts
  - mitre.attack.t1053.003   # Scheduled Task/Job: Cron
  - mitre.attack.t1548.001   # Abuse Elevation Control Mechanism: Setuid/Setgid
  - tlp.white
  - data_source.osquery
  - data_source.fim
category: endpoint
log_source:
  product: osquery
  service: file_events
detection:
  condition: |
    action IN ["CREATED", "DELETED", "UPDATED", "ATTRIBUTES_MODIFIED"]
    AND target_path STARTS_WITH_ANY [
      "/etc/passwd",
      "/etc/shadow",
      "/etc/gshadow",
      "/etc/sudoers",
      "/etc/sudoers.d/",
      "/etc/crontab",
      "/etc/cron.",
      "/etc/hosts",
      "/etc/ssh/sshd_config",
      "/etc/ld.so.conf",
      "/etc/ld.so.preload"
    ]
false_positives:
  - System package manager updates (apt, yum, dnf) touching /etc/passwd, /etc/shadow
  - Legitimate SSH configuration changes by ops team
  - cron package installation/update modifying /etc/crontab
  - SSO/LDAP integration tooling updating /etc/nsswitch.conf
playbook: tpl-credential-access
enabled: true
author: AiSOC
references:
  - https://attack.mitre.org/techniques/T1003/008/
  - https://attack.mitre.org/techniques/T1053/003/
created: '2026-05-10'
modified: '2026-05-10'
