# detections/endpoint/macos-launchagent-write.yaml
id: det-endpoint-031
name: macOS LaunchAgent or LaunchDaemon Plist Created
description: AiSOC v1 curated detection. Triggers on the endpoint signal described by 'macOS LaunchAgent
  or LaunchDaemon Plist Created'. Watch the 1 documented false-positive case before tuning.
version: 1.0.0
severity: high
tags:
- mitre.attack.t1543.001
- mitre.attack.t1543.004
- tlp.white
category: endpoint
log_source:
  product: macos
  service: endpointsecurity
detection:
  condition: |
    event_type == "ES_EVENT_TYPE_NOTIFY_CREATE"
    AND file_path STARTSWITH_ANY ["/Library/LaunchAgents/", "/Library/LaunchDaemons/", "/Users/Shared/"]
false_positives:
- Software installers from notarized publishers (allowlist signing-team-id)
playbook: tpl-persistence
enabled: true
author: AiSOC
created: '2026-05-03'
modified: '2026-05-03'
